# phase2_backup_and_recovery/systemd/db-dump.service
[Unit]
Description=Nightly DB dump
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/db_dump.sh
